# deploy/quadlet/qq-agent-app.container
[Unit]
Description=QQ Agent AI and graph-memory API

[Container]
Image=ghcr.io/vesperglow/qq-agent-app:latest
ContainerName=qq-agent-app
Network=qq-agent.network
EnvironmentFile=./app.env
PublishPort=127.0.0.1:8000:8000
AutoUpdate=registry
HealthCmd=python -c "import urllib.request; urllib.request.urlopen('http://127.0.0.1:8000/health/live', timeout=3)"
HealthInterval=15s
HealthRetries=5
HealthStartPeriod=20s
HealthTimeout=5s
HealthOnFailure=kill
Notify=healthy

[Service]
Restart=always
RestartSec=5s
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=default.target

